# .config/autostart/load-Xresources.desktop
[Desktop Entry]
Exec=xrdb -load ~/.Xresources
Name=load-Xresources
Type=Application
Version=1.0
