# .config/autostart/lxqt-config-monitor-autostart.desktop
[Desktop Entry]
Exec=lxqt-config-monitor -l
Name=lxqt-config-monitor-autostart
Type=Application
Version=1.0
