# .config/autostart/lxqt-config-touchpad-autostart.desktop
[Desktop Entry]
Exec=lxqt-config-input  --load-touchpad
Name=lxqt-config-touchpad-autostart
Type=Application
Version=1.0
