# .config/autostart/meteo-qt.desktop
[Desktop Entry]
Exec=/usr/bin/meteo-qt
Name=meteo-qt
Type=Application
Version=1.0
X-LXQt-Need-Tray=true
