# .config/autostart/wallpapers.desktop
[Desktop Entry]
Exec=/home/glaulher/.config/scripts/wallpapers.sh
Name=wallpapers
Type=Application
Version=1.0
