# .config/autostart/xbacklight.desktop
[Desktop Entry]
Exec=xbacklight -set 35
Name=xbacklight
Type=Application
Version=1.0
